# home/cuki/.config/systemd/user/xmonad.service
[Unit]
Description=Xmonad window manager

[Service]
ExecStart=/usr/bin/xmonad
Environment=DISPLAY=:0

[Install]
WantedBy=wm.target
