# systemd/user/startx.service
[Unit]
Description=StartX

[Service]
#User=bob
ExecStart=/usr/bin/startx

[Install]
WantedBy=graphical.target
